# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
redirect_stderr=true

[program:cron]
command=/usr/sbin/cron
autorestart=true
redirect_stderr=true
